# deploy/systemd/brunoosbrain-reflect-curate.timer
[Unit]
Description=brunoosbrain memory curation (daily-log distill + buffer drain + evict) — daily@08:00

[Timer]
OnCalendar=*-*-* 08:00:00 America/Sao_Paulo
Persistent=true
Unit=brunoosbrain-reflect-curate.service

[Install]
WantedBy=timers.target
