# deploy/systemd/brunoosbrain-vault-sync.service
[Unit]
Description=brunoosbrain vault git-sync (BrunOS/ ↔ GitHub)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=bruno
Group=bruno
WorkingDirectory=/home/bruno/BrunOS
Environment=TZ=America/Sao_Paulo
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/git-sync
TimeoutStartSec=120
LogsDirectory=brunoosbrain
StandardOutput=append:/var/log/brunoosbrain/vault-sync.log
StandardError=append:/var/log/brunoosbrain/vault-sync.log
