# deploy/systemd/brunoosbrain-slackbot.service
[Unit]
Description=brunoosbrain Slack chat bot (Socket Mode, single-instance)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=on-failure
RestartSec=10
User=bruno
Group=bruno
WorkingDirectory=/home/bruno/claude-second-brain
EnvironmentFile=/home/bruno/claude-second-brain/.claude/.env
Environment=TZ=America/Sao_Paulo
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/uv run python .claude/chat/bot.py
LogsDirectory=brunoosbrain
StandardOutput=append:/var/log/brunoosbrain/slackbot.log
StandardError=append:/var/log/brunoosbrain/slackbot.log

[Install]
WantedBy=multi-user.target
