# deploy/systemd/brunoosbrain-reflect-inbox.timer
[Unit]
Description=brunoosbrain inbox pass — hourly

[Timer]
OnCalendar=*-*-* 08..20:00:00
Persistent=true
Unit=brunoosbrain-reflect-inbox.service

[Install]
WantedBy=timers.target
