# deploy/systemd/brunoosbrain-weekly-review.service
[Unit]
Description=brunoosbrain weekly review draft
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=bruno
Group=bruno
WorkingDirectory=/home/bruno/claude-second-brain
EnvironmentFile=/home/bruno/claude-second-brain/.claude/.env
Environment=TZ=America/Sao_Paulo
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/uv run python .claude/skills/weekly-review/scripts/aggregate_week.py
TimeoutStartSec=600
LogsDirectory=brunoosbrain
StandardOutput=append:/var/log/brunoosbrain/weekly-review.log
StandardError=append:/var/log/brunoosbrain/weekly-review.log
