# deploy/systemd/brunoosbrain-heartbeat.service
[Unit]
Description=brunoosbrain heartbeat
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=bruno
Group=bruno
WorkingDirectory=/home/bruno/claude-second-brain
EnvironmentFile=/home/bruno/claude-second-brain/.claude/.env
Environment=TZ=America/Sao_Paulo
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/uv run python .claude/scripts/heartbeat.py
TimeoutStartSec=300
LogsDirectory=brunoosbrain
StandardOutput=append:/var/log/brunoosbrain/heartbeat.log
StandardError=append:/var/log/brunoosbrain/heartbeat.log
