# deploy/systemd/brunoosbrain-code-sync.service
[Unit]
Description=brunoosbrain code git-sync (claude-second-brain ↔ GitHub, pull-only)
After=network-online.target
Wants=network-online.target
OnFailure=brunoosbrain-alert@%n.service

[Service]
Type=oneshot
User=bruno
Group=bruno
WorkingDirectory=/home/bruno/claude-second-brain
Environment=TZ=America/Sao_Paulo
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/uv run python .claude/scripts/code_sync.py
TimeoutStartSec=120
LogsDirectory=brunoosbrain
StandardOutput=append:/var/log/brunoosbrain/code-sync.log
StandardError=append:/var/log/brunoosbrain/code-sync.log
