# deploy/systemd/linosbrain-dream.service
[Unit]
Description=linosbrain company-brain dreaming
After=network-online.target
Wants=network-online.target
OnFailure=linosbrain-alert@%n.service

[Service]
Type=oneshot
User=linos
Group=linos
WorkingDirectory=/home/linos/claude-second-brain
EnvironmentFile=/home/linos/claude-second-brain/.claude/.env
Environment=TZ=America/Sao_Paulo
Environment=COMPANY_BRAIN_PROFILE=linos
Environment=COMPANY_BRAIN_NAME=LinOS
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/uv run python .claude/scripts/company_brain_reflect.py dream --profile linos
TimeoutStartSec=600
LogsDirectory=linosbrain
StandardOutput=append:/var/log/linosbrain/dream.log
StandardError=append:/var/log/linosbrain/dream.log
