# deploy/systemd/linosbrain-reflect.timer
[Unit]
Description=linosbrain company-brain reflection - daily 09:15 BRT (after consumer)

[Timer]
OnCalendar=*-*-* 09:15 America/Sao_Paulo
Persistent=true
Unit=linosbrain-reflect.service

[Install]
WantedBy=timers.target
