# deploy/systemd/linosbrain-vault-sync.timer
[Unit]
Description=linosbrain vault sync — every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
Unit=linosbrain-vault-sync.service

[Install]
WantedBy=timers.target
