# deploy/launchd/com.bruno.brunos.chat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bruno.brunos.chat</string>

    <key>Disabled</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/brunobouwman/Documents/brunOS-brain</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/brunobouwman/.local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>.claude/chat/bot.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>America/Sao_Paulo</string>
        <key>PATH</key>
        <string>/Users/brunobouwman/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/brunobouwman/Library/Logs/com.bruno.brunos.chat.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/brunobouwman/Library/Logs/com.bruno.brunos.chat.log</string>
</dict>
</plist>
