# deploy/launchd/com.bruno.brunos.git-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bruno.brunos.git-sync</string>

    <key>Disabled</key>
    <false/>

    <key>RunAtLoad</key>
    <true/>

    <!-- Run via uv (holds Full Disk Access, granted for codex-watcher) so launchd
         can reach the ~/Documents-resident vault despite macOS TCC. WorkingDirectory
         is the repo (uv finds pyproject + .claude/.env there); vault_sync.py resolves
         the vault via shared.vault_path() and runs git there. A direct git exec on the
         vault hits "Operation not permitted"; going through uv inherits FDA to git. -->
    <key>WorkingDirectory</key>
    <string>/Users/brunobouwman/Documents/brunOS-brain</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/brunobouwman/.local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>.claude/scripts/vault_sync.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>America/Sao_Paulo</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/brunobouwman/Library/Logs/com.bruno.brunos.git-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/brunobouwman/Library/Logs/com.bruno.brunos.git-sync.log</string>

    <key>StartInterval</key>
    <integer>120</integer>
</dict>
</plist>
