# deploy/launchd/com.bruno.brunos.reflect-inbox.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bruno.brunos.reflect-inbox</string>

    <!-- Failover unit: ships disabled on the Mac (reflection/dream run on the
         VPS; MEMORY.md/playbook writes aren't dual-run safe). Operator enables
         on failover — see deploy/bin/install-mac-launchd.sh. -->
    <key>Disabled</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/brunobouwman/Documents/brunOS-brain</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/brunobouwman/.local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>.claude/scripts/memory_reflect.py</string>
        <string>--inbox-only</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>America/Sao_Paulo</string>
        <key>PATH</key>
        <string>/Users/brunobouwman/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/brunobouwman/Library/Logs/com.bruno.brunos.reflect-inbox.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/brunobouwman/Library/Logs/com.bruno.brunos.reflect-inbox.log</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>10</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>11</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>12</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>13</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>14</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>17</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>18</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>19</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>20</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>
</dict>
</plist>
